# docker-images/buggyapp/services.d/yc-agent/run
#!/usr/bin/with-contenv bash

cd /opt/workspace/yc/agent
./yc -c config.yaml -a buggyapp

